# server/crontab
* * * * * node /app/cron-job.js >> /app/cron.log 2>&1
